# ramepkg/rameplayer/wolfson.rules
SUBSYSTEM=="sound", ACTION=="add|change", ATTRS{id}=="RPiCirrus", \
	RUN+="/usr/libexec/wolfson/wolfson.sh"
