# ramepkg/rameplayer-utils/rameclock-led.initd
#!/sbin/openrc-run

name=rameclock-led
command=/usr/bin/$name
command_background=1
pidfile=/var/run/$name

depend() {
	need mosquitto
}
